# BangBangBandito/gallium-bandit.list
/usr
/usr/share
/usr/share/doc
/usr/share/doc/galliumos-artwork
/usr/share/doc/galliumos-artwork/changelog.gz
/usr/share/pixmaps
/usr/share/pixmaps/galliumos-solidblue.png
/usr/share/pixmaps/galliumos-trianglify-blue-1920.png
/usr/share/pixmaps/galliumos-trianglify-dkblue-1920.png
/usr/share/pixmaps/xubuntu-logo-menu.png
/usr/share/pixmaps/xubuntu-logo.png
/usr/share/pixmaps/xubuntu-logo.svg
/etc
/etc/apt
/etc/apt/preferences.d
/etc/apt/preferences.d/galliumos.pref
/etc/apt/sources.list.d
/etc/apt/sources.list.d/galliumos.list
/etc/apt/trusted.gpg.d
/etc/apt/trusted.gpg.d/galliumos.gpg
/etc/default
/etc/default/grub.d
/etc/default/grub.d/00-galliumos.cfg
/etc/default/grub.d/README
/etc/initramfs-tools
/etc/initramfs-tools/conf.d
/etc/initramfs-tools/conf.d/resume
/etc/modules-load.d
/etc/modules-load.d/cups-filters.conf
/etc/sysctl.d
/etc/sysctl.d/20-galliumos.conf
/etc/systemd
/etc/systemd/logind.conf
/etc/udev
/etc/udev/rules.d
/etc/udev/rules.d/99-hdmi.rules
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system-sleep
/lib/systemd/system-sleep/hdmi-change.sh
/lib/systemd/system-sleep/wpasupplicant
/lib/systemd/system/galliumos-base.service
/usr
/usr/bin
/usr/bin/change-brightness
/usr/bin/cpv
/usr/bin/firmware-util
/usr/bin/galliumos-repodist
/usr/bin/hdmi-change
/usr/bin/init-galliumos-opts
/usr/share
/usr/share/applications
/usr/share/applications/galliumos-irc.desktop
/usr/share/doc
/usr/share/doc/galliumos-base
/usr/share/doc/galliumos-base/changelog.gz
/usr/share/galliumos-repodist
/usr/share/galliumos-repodist/dists
/usr/share/galliumos-repodist/dists/bismuth-devel
/usr/share/galliumos-repodist/dists/bismuth-devel/preferences.d
/usr/share/galliumos-repodist/dists/bismuth-devel/preferences.d/bismuth-devel.pref
/usr/share/galliumos-repodist/dists/bismuth-devel/sources.list.d
/usr/share/galliumos-repodist/dists/bismuth-devel/sources.list.d/bismuth-devel.list
/usr/share/galliumos-repodist/dists/bismuth-prerelease
/usr/share/galliumos-repodist/dists/bismuth-prerelease/preferences.d
/usr/share/galliumos-repodist/dists/bismuth-prerelease/preferences.d/bismuth-prerelease.pref
/usr/share/galliumos-repodist/dists/bismuth-prerelease/sources.list.d
/usr/share/galliumos-repodist/dists/bismuth-prerelease/sources.list.d/bismuth-prerelease.list
/usr/share/galliumos-repodist/dists/bismuth-testing
/usr/share/galliumos-repodist/dists/bismuth-testing/preferences.d
/usr/share/galliumos-repodist/dists/bismuth-testing/preferences.d/bismuth-testing.pref
/usr/share/galliumos-repodist/dists/bismuth-testing/sources.list.d
/usr/share/galliumos-repodist/dists/bismuth-testing/sources.list.d/bismuth-testing.list
/etc
/etc/acpi
/etc/acpi/events
/etc/acpi/events/jack
/etc/acpi/jack.sh
/etc/modprobe.d
/etc/modprobe.d/alsa-chtrt5650.conf
/etc/modprobe.d/iwlwifi-opts.conf
/etc/pulse
/etc/pulse/default.pa
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system-sleep
/lib/systemd/system-sleep/pulseaudio
/lib/systemd/system-sleep/wifi.sh
/lib/systemd/system/galliumos-braswell.service
/usr
/usr/bin
/usr/bin/init-braswell-settings
/usr/share
/usr/share/X11
/usr/share/X11/xorg.conf.d
/usr/share/X11/xorg.conf.d/20-intel.conf
/usr/share/alsa
/usr/share/alsa/ucm
/usr/share/alsa/ucm/chtmax98090
/usr/share/alsa/ucm/chtmax98090/HiFi.conf
/usr/share/alsa/ucm/chtmax98090/chtmax98090.conf
/usr/share/alsa/ucm/chtrt5650
/usr/share/alsa/ucm/chtrt5650/HiFi.conf
/usr/share/alsa/ucm/chtrt5650/chtrt5650.conf
/usr/share/doc
/usr/share/doc/galliumos-braswell
/usr/share/doc/galliumos-braswell/changelog.gz
/usr/share/pulseaudio
/usr/share/pulseaudio/alsa-mixer
/usr/share/pulseaudio/alsa-mixer/paths
/usr/share/pulseaudio/alsa-mixer/paths/analog-output-headphones.conf
/usr/share/pulseaudio/alsa-mixer/paths/analog-output-speaker.conf
/usr
/usr/share
/usr/share/doc
/usr/share/doc/galliumos-core
/usr/share/doc/galliumos-core/changelog.gz
/usr/share/doc/galliumos-core/copyright
/usr
/usr/share
/usr/share/doc
/usr/share/doc/galliumos-default-settings
/usr/share/doc/galliumos-default-settings/changelog.gz
/usr/share/doc/galliumos-default-settings/copyright
/usr
/usr/share
/usr/share/applications
/usr/share/applications/display-im6.desktop
/usr/share/applications/display-im6.q16.desktop
/usr/share/doc
/usr/share/doc/galliumos-desktop
/usr/share/doc/galliumos-desktop/changelog.gz
/usr/share/doc/galliumos-desktop/copyright
/etc
/etc/skel
/etc/skel/.config
/etc/skel/.config/xfce4
/etc/skel/.config/xfce4/xfconf
/etc/skel/.config/xfce4/xfconf/xfce-perchannel-xml
/etc/skel/.config/xfce4/xfconf/xfce-perchannel-xml/keyboard-layout.xml
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/galliumos-laptop.service
/usr
/usr/bin
/usr/bin/init-laptop-settings
/usr/share
/usr/share/doc
/usr/share/doc/galliumos-laptop
/usr/share/doc/galliumos-laptop/changelog.gz
/usr
/usr/share
/usr/share/doc
/usr/share/doc/galliumos-wallpapers
/usr/share/doc/galliumos-wallpapers/changelog.gz
/usr/share/xfce4
/usr/share/xfce4/backdrops
/usr/share/xfce4/backdrops/README-galliumos.md
/usr/share/xfce4/backdrops/galliumos-1.0-architectural.jpg
/usr/share/xfce4/backdrops/galliumos-1.0-solidblue.png
/usr/share/xfce4/backdrops/galliumos-2.0-trianglify-blue.jpg
/usr/share/xfce4/backdrops/galliumos-2.1-architectural.jpg
/usr/share/xfce4/backdrops/galliumos-2.1-turbines.jpg
/usr/share/xfce4/backdrops/galliumos-2.1-underwater.jpg
/usr/share/xfce4/backdrops/galliumos-3.0-architectural.jpg
/usr/share/xfce4/backdrops/galliumos-3.0-graymist.png
/etc
/etc/apt
/etc/apt/preferences.d
/etc/apt/preferences.d/galliumos.pref
/etc/apt/sources.list.d
/etc/apt/sources.list.d/galliumos.list
/etc/apt/trusted.gpg.d
/etc/apt/trusted.gpg.d/galliumos.gpg
/etc/default
/etc/default/grub.d
/etc/default/grub.d/00-galliumos.cfg
/etc/default/grub.d/README
/etc/initramfs-tools
/etc/initramfs-tools/conf.d
/etc/initramfs-tools/conf.d/resume
/etc/modules-load.d
/etc/modules-load.d/cups-filters.conf
/etc/sysctl.d
/etc/sysctl.d/20-galliumos.conf
/etc/systemd
/etc/systemd/logind.conf
/etc/udev
/etc/udev/rules.d
/etc/udev/rules.d/99-hdmi.rules
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system-sleep
/lib/systemd/system-sleep/hdmi-change.sh
/lib/systemd/system-sleep/wpasupplicant
/lib/systemd/system/galliumos-base.service
/usr
/usr/bin
/usr/bin/change-brightness
/usr/bin/cpv
/usr/bin/firmware-util
/usr/bin/galliumos-repodist
/usr/bin/hdmi-change
/usr/bin/init-galliumos-opts
/usr/share
/usr/share/applications
/usr/share/applications/galliumos-irc.desktop
/usr/share/doc
/usr/share/doc/galliumos-base
/usr/share/doc/galliumos-base/changelog.gz
/usr/share/galliumos-repodist
/usr/share/galliumos-repodist/dists
/usr/share/galliumos-repodist/dists/bismuth-devel
/usr/share/galliumos-repodist/dists/bismuth-devel/preferences.d
/usr/share/galliumos-repodist/dists/bismuth-devel/preferences.d/bismuth-devel.pref
/usr/share/galliumos-repodist/dists/bismuth-devel/sources.list.d
/usr/share/galliumos-repodist/dists/bismuth-devel/sources.list.d/bismuth-devel.list
/usr/share/galliumos-repodist/dists/bismuth-prerelease
/usr/share/galliumos-repodist/dists/bismuth-prerelease/preferences.d
/usr/share/galliumos-repodist/dists/bismuth-prerelease/preferences.d/bismuth-prerelease.pref
/usr/share/galliumos-repodist/dists/bismuth-prerelease/sources.list.d
/usr/share/galliumos-repodist/dists/bismuth-prerelease/sources.list.d/bismuth-prerelease.list
/usr/share/galliumos-repodist/dists/bismuth-testing
/usr/share/galliumos-repodist/dists/bismuth-testing/preferences.d
/usr/share/galliumos-repodist/dists/bismuth-testing/preferences.d/bismuth-testing.pref
/usr/share/galliumos-repodist/dists/bismuth-testing/sources.list.d
/usr/share/galliumos-repodist/dists/bismuth-testing/sources.list.d/bismuth-testing.list
/etc
/etc/acpi
/etc/acpi/events
/etc/acpi/events/jack
/etc/acpi/jack.sh
/etc/modprobe.d
/etc/modprobe.d/alsa-chtrt5650.conf
/etc/modprobe.d/iwlwifi-opts.conf
/etc/pulse
/etc/pulse/default.pa
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system-sleep
/lib/systemd/system-sleep/pulseaudio
/lib/systemd/system-sleep/wifi.sh
/lib/systemd/system/galliumos-braswell.service
/usr
/usr/bin
/usr/bin/init-braswell-settings
/usr/share
/usr/share/X11
/usr/share/X11/xorg.conf.d
/usr/share/X11/xorg.conf.d/20-intel.conf
/usr/share/alsa
/usr/share/alsa/ucm
/usr/share/alsa/ucm/chtmax98090
/usr/share/alsa/ucm/chtmax98090/HiFi.conf
/usr/share/alsa/ucm/chtmax98090/chtmax98090.conf
/usr/share/alsa/ucm/chtrt5650
/usr/share/alsa/ucm/chtrt5650/HiFi.conf
/usr/share/alsa/ucm/chtrt5650/chtrt5650.conf
/usr/share/doc
/usr/share/doc/galliumos-braswell
/usr/share/doc/galliumos-braswell/changelog.gz
/usr/share/pulseaudio
/usr/share/pulseaudio/alsa-mixer
/usr/share/pulseaudio/alsa-mixer/paths
/usr/share/pulseaudio/alsa-mixer/paths/analog-output-headphones.conf
/usr/share/pulseaudio/alsa-mixer/paths/analog-output-speaker.conf
/usr
/usr/share
/usr/share/doc
/usr/share/doc/galliumos-core
/usr/share/doc/galliumos-core/changelog.gz
/usr/share/doc/galliumos-core/copyright
/etc
/etc/skel
/etc/skel/.config
/etc/skel/.config/chromium
/etc/skel/.config/chromium/Default
/etc/skel/.config/chromium/Default/Preferences
/etc/skel/.config/vlc
/etc/skel/.config/vlc/vlc-qt-interface.conf
/etc/skel/.config/vlc/vlcrc
/etc/skel/.config/xfce4
/etc/skel/.config/xfce4/helpers.rc
/etc/skel/.config/xfce4/panel
/etc/skel/.config/xfce4/panel/launcher-1
/etc/skel/.config/xfce4/panel/launcher-1/14449738521.desktop
/etc/skel/.config/xfce4/panel/launcher-12
/etc/skel/.config/xfce4/panel/launcher-12/14449738943.desktop
/etc/skel/.config/xfce4/panel/launcher-2
/etc/skel/.config/xfce4/panel/launcher-2/14449738812.desktop
/etc/skel/.config/xfce4/panel/whiskermenu-7.rc
/etc/skel/.config/xfce4/terminal
/etc/skel/.config/xfce4/terminal/terminalrc
/etc/skel/.config/xfce4/xfconf
/etc/skel/.config/xfce4/xfconf/xfce-perchannel-xml
/etc/skel/.config/xfce4/xfconf/xfce-perchannel-xml/keyboards.xml
/etc/skel/.config/xfce4/xfconf/xfce-perchannel-xml/thunar.xml
/etc/skel/.config/xfce4/xfconf/xfce-perchannel-xml/xfce4-desktop.xml
/etc/skel/.config/xfce4/xfconf/xfce-perchannel-xml/xfce4-keyboard-shortcuts.xml
/etc/skel/.config/xfce4/xfconf/xfce-perchannel-xml/xfce4-mixer.xml
/etc/skel/.config/xfce4/xfconf/xfce-perchannel-xml/xfce4-notifyd.xml
/etc/skel/.config/xfce4/xfconf/xfce-perchannel-xml/xfce4-panel.xml
/etc/skel/.config/xfce4/xfconf/xfce-perchannel-xml/xfce4-power-manager.xml
/etc/skel/.config/xfce4/xfconf/xfce-perchannel-xml/xfce4-session.xml
/etc/skel/.config/xfce4/xfconf/xfce-perchannel-xml/xfce4-settings-manager.xml
/etc/skel/.config/xfce4/xfconf/xfce-perchannel-xml/xfwm4.xml
/etc/skel/.config/xfce4/xfconf/xfce-perchannel-xml/xsettings.xml
/etc/skel/.local
/etc/skel/.local/share
/etc/skel/.local/share/applications
/etc/skel/.local/share/applications/mimeapps.list
/etc/skel/.xscreensaver
/usr
/usr/share
/usr/share/doc
/usr/share/doc/galliumos-default-settings
/usr/share/doc/galliumos-default-settings/changelog.gz
/usr/share/doc/galliumos-default-settings/copyright
/usr
/usr/share
/usr/share/applications
/usr/share/applications/display-im6.desktop
/usr/share/applications/display-im6.q16.desktop
/usr/share/doc
/usr/share/doc/galliumos-desktop
/usr/share/doc/galliumos-desktop/changelog.gz
/usr/share/doc/galliumos-desktop/copyright
/etc
/etc/skel
/etc/skel/.config
/etc/skel/.config/xfce4
/etc/skel/.config/xfce4/xfconf
/etc/skel/.config/xfce4/xfconf/xfce-perchannel-xml
/etc/skel/.config/xfce4/xfconf/xfce-perchannel-xml/keyboard-layout.xml
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/galliumos-laptop.service
/usr
/usr/bin
/usr/bin/init-laptop-settings
/usr/share
/usr/share/doc
/usr/share/doc/galliumos-laptop
/usr/share/doc/galliumos-laptop/changelog.gz
/etc
/etc/skel
/etc/skel/.config
/etc/skel/.config/xfce4
/etc/skel/.config/xfce4/xfconf
/etc/skel/.config/xfce4/xfconf/xfce-perchannel-xml
/etc/skel/.config/xfce4/xfconf/xfce-perchannel-xml/keyboard-layout.xml
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/galliumos-laptop.service
/usr
/usr/bin
/usr/bin/init-laptop-settings
/usr/share
/usr/share/doc
/usr/share/doc/galliumos-laptop
/usr/share/doc/galliumos-laptop/changelog.gz
/usr
/usr/share
/usr/share/doc
/usr/share/doc/galliumos-wallpapers
/usr/share/doc/galliumos-wallpapers/changelog.gz
/usr/share/xfce4
/usr/share/xfce4/backdrops
/usr/share/xfce4/backdrops/README-galliumos.md
/usr/share/xfce4/backdrops/galliumos-1.0-architectural.jpg
/usr/share/xfce4/backdrops/galliumos-1.0-solidblue.png
/usr/share/xfce4/backdrops/galliumos-2.0-trianglify-blue.jpg
/usr/share/xfce4/backdrops/galliumos-2.1-architectural.jpg
/usr/share/xfce4/backdrops/galliumos-2.1-turbines.jpg
/usr/share/xfce4/backdrops/galliumos-2.1-underwater.jpg
/usr/share/xfce4/backdrops/galliumos-3.0-architectural.jpg
/usr/share/xfce4/backdrops/galliumos-3.0-graymist.png
